# local/.config/systemd/user/kubeunlock.timer
[Unit]
Description=Unlock kubeconfig at 5pm on workdays

[Timer]
OnCalendar=Mon..Fri 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
